# packaging/propertyquarry-release-control-v2/systemd/propertyquarry-release-control-v2@.service
[Unit]
Description=PropertyQuarry release-control v2 isolated server broker
Documentation=man:systemd.service(5)
Requires=propertyquarry-release-control-v2.socket
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3
RefuseManualStart=yes
CollectMode=inactive-or-failed

[Service]
Type=exec
ExitType=cgroup
User=propertyquarry-release-control
Group=propertyquarry-release-control
ExecStart=/usr/libexec/propertyquarry-release-control/propertyquarry-release-supervisor-v2 --server-broker --config /etc/propertyquarry-release-control/controller-v2.json --socket-activation
WorkingDirectory=/var/lib/propertyquarry-release-control-v2
StandardInput=socket
StandardOutput=journal
StandardError=journal
SyslogIdentifier=propertyquarry-release-control-v2
LogLevelMax=notice
LogRateLimitIntervalSec=30
LogRateLimitBurst=50
UMask=0077
StateDirectory=propertyquarry-release-control-v2
StateDirectoryMode=0700
LoadCredentialEncrypted=request-authority-client:/etc/credstore.encrypted/propertyquarry-release-control-v2-request-authority-client
LoadCredentialEncrypted=lifecycle-cas-client:/etc/credstore.encrypted/propertyquarry-release-control-v2-lifecycle-cas-client
LoadCredentialEncrypted=evidence-store-client:/etc/credstore.encrypted/propertyquarry-release-control-v2-evidence-store-client
LoadCredentialEncrypted=resource-mediator-client:/etc/credstore.encrypted/propertyquarry-release-control-v2-resource-mediator-client
RuntimeMaxSec=9900
TimeoutStopSec=15
KillMode=control-group
SendSIGKILL=yes
FinalKillSignal=SIGKILL
Delegate=no
OOMPolicy=kill
Restart=no
TasksMax=256
LimitNOFILE=4096
LimitCORE=0
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
SecureBits=noroot-locked
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes
KeyringMode=private
DevicePolicy=closed
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@clock @keyring @mount @obsolete @privileged @raw-io @reboot @resources @swap
UnsetEnvironment=ACTIONS_ID_TOKEN_REQUEST_TOKEN ACTIONS_ID_TOKEN_REQUEST_URL PROPERTYQUARRY_OIDC_TOKEN_FD GITHUB_ACTION GITHUB_ACTOR GITHUB_ENV GITHUB_EVENT_PATH GITHUB_HEAD_REF GITHUB_JOB GITHUB_PATH GITHUB_REF GITHUB_REPOSITORY GITHUB_RUN_ATTEMPT GITHUB_RUN_ID GITHUB_SHA GITHUB_TOKEN GITHUB_WORKFLOW GITHUB_WORKFLOW_REF GITHUB_WORKFLOW_SHA GITHUB_WORKSPACE RUNNER_TEMP RUNNER_TOOL_CACHE HTTP_PROXY HTTPS_PROXY ALL_PROXY NO_PROXY http_proxy https_proxy all_proxy no_proxy SSL_CERT_FILE SSL_CERT_DIR CURL_CA_BUNDLE GIT_SSL_CAINFO LD_PRELOAD LD_LIBRARY_PATH LD_AUDIT PYTHONPATH PYTHONHOME RUST_LOG RUST_BACKTRACE GODEBUG
