# packaging/propertyquarry-release-control-v2/systemd/propertyquarry-release-watchdog-v2.service
[Unit]
Description=PropertyQuarry release-control v2 reconciliation watchdog
Documentation=man:systemd.service(5)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5
CollectMode=inactive-or-failed

[Service]
Type=notify
ExitType=cgroup
NotifyAccess=main
User=propertyquarry-release-watchdog
Group=propertyquarry-release-watchdog
SupplementaryGroups=propertyquarry-release-control
ExecStart=/usr/libexec/propertyquarry-release-control/propertyquarry-release-watchdog-v2 --config /etc/propertyquarry-release-control/watchdog-v2.json
WorkingDirectory=/var/lib/propertyquarry-release-watchdog-v2
StandardInput=null
StandardOutput=journal
StandardError=journal
SyslogIdentifier=propertyquarry-release-watchdog-v2
LogLevelMax=notice
LogRateLimitIntervalSec=30
LogRateLimitBurst=50
UMask=0077
StateDirectory=propertyquarry-release-watchdog-v2
StateDirectoryMode=0700
RuntimeDirectory=propertyquarry-release-watchdog-v2
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
LoadCredentialEncrypted=watchdog-takeover-client:/etc/credstore.encrypted/propertyquarry-release-control-v2-watchdog-takeover-client
LoadCredentialEncrypted=resource-recovery-client:/etc/credstore.encrypted/propertyquarry-release-control-v2-resource-recovery-client
WatchdogSec=30
TimeoutStartSec=60
TimeoutStopSec=15
KillMode=control-group
SendSIGKILL=yes
FinalKillSignal=SIGKILL
Delegate=no
OOMPolicy=kill
Restart=always
RestartSec=5
TasksMax=128
LimitNOFILE=2048
LimitCORE=0
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
SecureBits=noroot-locked
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes
KeyringMode=private
DevicePolicy=closed
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@clock @keyring @mount @obsolete @privileged @raw-io @reboot @resources @swap
UnsetEnvironment=ACTIONS_ID_TOKEN_REQUEST_TOKEN ACTIONS_ID_TOKEN_REQUEST_URL PROPERTYQUARRY_OIDC_TOKEN_FD GITHUB_ACTION GITHUB_ACTOR GITHUB_ENV GITHUB_EVENT_PATH GITHUB_HEAD_REF GITHUB_JOB GITHUB_PATH GITHUB_REF GITHUB_REPOSITORY GITHUB_RUN_ATTEMPT GITHUB_RUN_ID GITHUB_SHA GITHUB_TOKEN GITHUB_WORKFLOW GITHUB_WORKFLOW_REF GITHUB_WORKFLOW_SHA GITHUB_WORKSPACE RUNNER_TEMP RUNNER_TOOL_CACHE HTTP_PROXY HTTPS_PROXY ALL_PROXY NO_PROXY http_proxy https_proxy all_proxy no_proxy SSL_CERT_FILE SSL_CERT_DIR CURL_CA_BUNDLE GIT_SSL_CAINFO LD_PRELOAD LD_LIBRARY_PATH LD_AUDIT PYTHONPATH PYTHONHOME RUST_LOG RUST_BACKTRACE GODEBUG

[Install]
WantedBy=multi-user.target
